# dotfiles/systemd/user/syncthing.service
[Unit]
Description=SyncThingy Daemon

[Service]
ExecStart=flatpak run com.github.zocker_160.SyncThingy

[Install]
WantedBy=default.target
