# dotfiles/brian-pop-os-pc-specialist/autostart/devilspie2.desktop
[Desktop Entry]
Type=Application
Name=Devilspie2 autostart
Exec=devilspie2
Comment=Start Devilspie2 on desktop start
