# dotfiles/brian-pop-os-pc-specialist/autostart/slack.desktop
[Desktop Entry]
Type=Application
Name=Slack autostart
Exec=flatpak run com.slack.Slack
Comment=Start Slack on desktop start
