# dotfiles/brian-pop-os-pc-specialist/autostart/wezterm.desktop
[Desktop Entry]
Type=Application
Name=Wezterm autostart
Exec=flatpak run org.wezfurlong.wezterm
Comment=Start Wezterm on desktop start
