# dotfiles/brians-infinitybookpro14-popos/autostart/com.dropbox.Client.desktop
[Desktop Entry]
Type=Application
Name=com.dropbox.Client
Exec=flatpak run com.dropbox.Client
X-Flatpak=com.dropbox.Client
