# dotfiles/brians-infinitybookpro14-popos/autostart/discord.desktop
[Desktop Entry]
Type=Application
Name=Discord autostart
Exec=flatpak run com.discordapp.Discord
Comment=Start Discord on desktop start
