# dotfiles/brians-infinitybookpro14-popos/autostart/firefox.desktop
[Desktop Entry]
Type=Application
Name=Firefox autostart
Exec=firefox
Comment=Start Firefox on desktop start
