# dotfiles/brians-infinitybookpro14-popos/autostart/thunderbird.desktop
[Desktop Entry]
Type=Application
Name=Thunderbird autostart
Exec=flatpak run org.mozilla.Thunderbird
Comment=Start Thunderbird on desktop start
