# dotfiles/autostart/firefox.desktop
[Desktop Entry]
Type=Application
Name=Firefox
Exec=firefox
Comment=Start Firefox on desktop start
