# dotfiles/brians-xps13-popos/autostart/yakuake.desktop
[Desktop Entry]
Type=Application
Name=Yakuake autostart
Exec=yakuake
Comment=Start Yakuake on desktop start
